# server-ce/runit/git-overleaf/run
#!/bin/bash

NODE_PARAMS=""
if [ "$DEBUG_NODE" == "true" ]; then
    echo "running debug - git"
    NODE_PARAMS="--inspect=0.0.0.0:30990"
fi

source /etc/overleaf/env.sh
export LISTEN_ADDRESS=127.0.0.1

exec /sbin/setuser www-data /usr/bin/node $NODE_PARAMS /overleaf/services/git/app.js >> /var/log/overleaf/git.log 2>&1
